# supervisorctl笔记/admin-server-main.conf
[program:admin-server-main]
directory = /data/admin-server-main
command = /data/admin-server-main/bin/admin-server-main -app-conf=/data/admin-server-main/appconf/app.conf
autostart = true
startsecs = 5
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisord/admin-server-main.log

